# ansible/roles/application/templates/3dfigurine-trellis.service.j2
[Unit]
Description=3D Figurine Lab - TRELLIS.2 Engine
After=docker.service
Requires=docker.service
PartOf=3dfigurine.target

[Service]
Type=simple
Restart={{ docker_restart_policy }}
RestartSec=10
User=azureuser
WorkingDirectory=/opt/3dfigurine

Environment="HOST_UID=1000"
Environment="HOST_GID=1000"
Environment="DOCKER_HOST=unix:///var/run/docker.sock"

ExecStartPre=/usr/bin/docker-compose pull 3dfigurine-trellis
ExecStart=/usr/bin/docker-compose up 3dfigurine-trellis
ExecStop=/usr/bin/docker-compose down 3dfigurine-trellis

StandardOutput=journal
StandardError=journal
SyslogIdentifier=3dfigurine-trellis

TimeoutStartSec=0
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
